# supervisord.conf
[supervisord]
nodaemon=true

[program:prosody]
command=/usr/bin/prosody
user=prosody

